# 红书运营/routines/launchd/com.hongshu.daily-topic-april.example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.hongshu.daily-topic-april</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-c</string>
		<string>export PATH="/opt/homebrew/bin:/usr/bin:/bin"; export TZ=Asia/Shanghai; exec "/Users/zhaoziqi/Documents/红书运营/routines/daily_run_topic_april.sh"</string>
	</array>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin</string>
		<key>TZ</key>
		<string>Asia/Shanghai</string>
	</dict>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>7</integer>
		<key>Minute</key>
		<integer>0</integer>
	</dict>
	<key>StandardOutPath</key>
	<string>/Users/zhaoziqi/Documents/红书运营/routines/logs/launchd-stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/zhaoziqi/Documents/红书运营/routines/logs/launchd-stderr.log</string>
	<key>WorkingDirectory</key>
	<string>/Users/zhaoziqi/Documents/红书运营</string>
</dict>
</plist>
